# supervisord.conf
[supervisord]
nodaemon=true

[program:spark_master]
command=/usr/local/spark/start-spark.sh
